# conf/supervisor/conf/conf.d/sshd.conf
[program:openssh]
command=/usr/sbin/sshd -D