# conf/supervisor/conf/conf.d/zookeeper.conf
[program:zookeeper]
command=/usr/local/zookeeper/bin/zkServer.sh start-foreground